# src/background_jobs/supervisord.conf
[inet_http_server]
port=0.0.0.0:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://0.0.0.0:9001

[supervisord]

[program:remove_expired_users]
environment = PYTHONUNBUFFERED=1
user=root
directory=/app/src/background_jobs/
command=/usr/local/bin/python /app/src/background_jobs/remove_expired_user_accounts.py
autostart=true
autorestart=true
stderr_logfile=/var/log/remove_expired_user_accounts/remove_expired_user_accounts.err.log
stdout_logfile=/var/log/remove_expired_user_accounts/remove_expired_user_accounts.out.log